# mission-control/mission-control-web.service
[Unit]
Description=Mission Control Web Dashboard
After=network.target mission-control-watcher.service

[Service]
Type=simple
WorkingDirectory=/home/andrew/claudeclaw/mission-control
ExecStart=/home/andrew/.nvm/versions/node/v20.20.0/bin/npx next start -p 3000 -H 0.0.0.0
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=0.0.0.0
Environment=NODE_OPTIONS=--no-network-family-autoselection

[Install]
WantedBy=default.target
